# system/postmap@.service
[Unit]
Description=Run postmap on %f

[Service]
Type=oneshot
ExecStart=postmap %f
